# install/iptv_cloud/service/service_slave.service.in
[Unit]
Description=@PROJECT_SUMMARY@
After=redis.service

[Service]
Type=forking
PermissionsStartOnly=true
PIDFile=@PIDFILE_PATH@
ExecStartPre=-/bin/mkdir @RUN_DIR_PATH@
ExecStartPre=-/bin/touch @PIDFILE_PATH@
ExecStartPre=-/bin/chown -R @USER_NAME@:@USER_GROUP@ @RUN_DIR_PATH@ @PIDFILE_PATH@
ExecStart=@EXECUTABLE_PATH@ --daemon
ExecStop=@EXECUTABLE_PATH@ --stop
ExecStopPost=-/bin/rm -f @PIDFILE_PATH@
User=@USER_NAME@
Group=@USER_GROUP@

[Install]
WantedBy=multi-user.target
